# deployment/systemd/ad-detection-web.service
[Unit]
Description=Ad Detection Web Interface (Cluster Leader)
After=network.target
Wants=network.target
ConditionPathExists=/var/lib/ad-detection/cluster/is_leader
StartLimitIntervalSec=0

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/opt/ad-detection
Environment="PYTHONPATH=/opt/ad-detection/packages/edge-device/src"
ExecStart=/opt/ad-detection/venv/bin/uvicorn web_interface.app:app --host 0.0.0.0 --port 8080 --workers 2
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ad-detection-web

# Watchdog configuration
WatchdogSec=30

# Resource limits
MemoryLimit=256M
CPUQuota=50%

# Security hardening
PrivateTmp=true
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/var/lib/ad-detection

[Install]
WantedBy=multi-user.target
